# outreach-agent/deploy/systemd/moltbook-outreach-heartbeat.service
[Unit]
Description=Moltbook outreach heartbeat
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__REMOTE_USER__
Group=__REMOTE_USER__
WorkingDirectory=__PACKAGE_DIR__
TimeoutStartSec=30min
TimeoutStopSec=20s
KillMode=control-group
Environment=HOME=/home/__REMOTE_USER__
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=MOLTBOOK_AGENT_ID=__AGENT_ID__
Environment=OUTREACH_AGENT_VENUE=moltbook
Environment=MOLTBOOK_STATE_PATH=__RUNTIME_DIR__/state.json
EnvironmentFile=__ENV_FILE__
ExecStartPre=/usr/bin/mkdir -p __RUNTIME_DIR__
ExecStart=/usr/bin/env bash -lc 'exec 9>"__LOCK_FILE__"; if ! flock -n 9; then echo "Skipping heartbeat because another run is active."; exit 0; fi; exec node dist/src/index.js heartbeat'
